# rag-pipeline/systemd/rag-ingest.timer
[Unit]
Description=Run Obsidian Wiki RAG Ingestion hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
